# mpkind.service
[Unit]
Description=MPKind Next.js Application
After=network.target

[Service]
Type=simple
User=nic0
WorkingDirectory=/home/nic0/programs/mpkind
ExecStart=/usr/bin/npm start
Restart=always
Environment=NODE_ENV=production
Environment=PORT=8600
Environment=HTTP_USER=admin
Environment=HTTP_PASSWORD=your_secure_password_here

[Install]
WantedBy=multi-user.target 